# scripts/server/templates/celery-beat-instance.service.template
[Unit]
Description=Celery beat for docketworks __INSTANCE__
After=network.target redis-server.service

[Service]
User=__INSTANCE_USER__
Group=www-data
WorkingDirectory=/opt/docketworks/instances/__INSTANCE__
ExecStart=/opt/docketworks/.venv/bin/celery -A docketworks beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
EnvironmentFile=/opt/docketworks/instances/__INSTANCE__/.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
