# scripts/server/templates/scheduler-instance.service.template
[Unit]
Description=Scheduler instance for docketworks __INSTANCE__
After=network.target

[Service]
User=__INSTANCE_USER__
Group=www-data
WorkingDirectory=/opt/docketworks/instances/__INSTANCE__
ExecStart=/opt/docketworks/.venv/bin/python manage.py run_scheduler
EnvironmentFile=/opt/docketworks/instances/__INSTANCE__/.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
